# etc/systemd/system/mail-cron.service
[Unit]
Description=Mail cron

[Service]
Type=simple
ExecStart=/usr/bin/php /var/www/mail/cron.php
